# docker/railway/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; Caddy validates the Jaban Universe JWT and proxies to Paperclip on loopback.
[program:caddy]
command=/usr/bin/caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; Paperclip binds to loopback only — required for local_trusted deployment mode.
[program:paperclip]
command=node --import ./server/node_modules/tsx/dist/loader.mjs server/dist/index.js
directory=/app
user=node
environment=PORT=3000,HOST=127.0.0.1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
