# paperclip-dpo-service/ai.whitestag.paperclip-dpo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.whitestag.paperclip-dpo</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BIN__</string>
    <string>__INSTALL_DIR__/dist/index.js</string>
  </array>
  <key>WorkingDirectory</key><string>__INSTALL_DIR__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>DPO_SHARED_KEY</key><string>__SHARED_KEY__</string>
    <key>DPO_MAPPING_DB</key><string>/var/paperclip/dpo/mappings.db</string>
    <key>DPO_AUDIT_DIR</key><string>/var/paperclip/dpo/audit</string>
    <key>DPO_CLASSIFIER_URL</key><string>http://localhost:1234</string>
    <key>DPO_CLASSIFIER_MODEL</key><string>google/gemma-4-26b-a4b</string>
  </dict>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>30</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/var/log/paperclip-dpo/out.log</string>
  <key>StandardErrorPath</key><string>/var/log/paperclip-dpo/err.log</string>
</dict>
</plist>
