# docker/quadlet/aiteamcorp.container
[Unit]
Description=Paperclip AI Agent Orchestrator
Requires=paperclip-db.service
After=paperclip-db.service

[Container]
Image=paperclip-local
ContainerName=paperclip
Pod=aiteamcorp.pod
Volume=%h/.local/share/paperclip:/paperclip:Z
Environment=HOST=0.0.0.0
Environment=AITEAMCORP_HOME=/aiteamcorp
Environment=AITEAMCORP_DEPLOYMENT_MODE=authenticated
Environment=AITEAMCORP_DEPLOYMENT_EXPOSURE=private
Environment=AITEAMCORP_PUBLIC_URL=http://localhost:3100
EnvironmentFile=%h/.config/containers/systemd/paperclip.env

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
