# docker/quadlet/aiteamcorp.container
[Unit]
Description=AiTeamCorp AI Agent Orchestrator
Requires=aiteamcorp-db.service
After=aiteamcorp-db.service

[Container]
Image=aiteamcorp-local
ContainerName=aiteamcorp
Pod=aiteamcorp.pod
Volume=%h/.local/share/aiteamcorp:/aiteamcorp:Z
Environment=HOST=0.0.0.0
Environment=AITEAMCORP_HOME=/aiteamcorp
Environment=AITEAMCORP_DEPLOYMENT_MODE=authenticated
Environment=AITEAMCORP_DEPLOYMENT_EXPOSURE=private
Environment=AITEAMCORP_PUBLIC_URL=http://localhost:3100
EnvironmentFile=%h/.config/containers/systemd/aiteamcorp.env

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
