# docker/quadlet/paperclip.container
[Unit]
Description=Paperclip AI Agent Orchestrator
Requires=paperclip-db.service
After=paperclip-db.service

[Container]
Image=paperclip-local
ContainerName=paperclip
Pod=paperclip.pod
Volume=%h/.local/share/paperclip:/paperclip:Z
Environment=HOST=0.0.0.0
Environment=PAPERCLIP_HOME=/paperclip
Environment=PAPERCLIP_DEPLOYMENT_MODE=authenticated
Environment=PAPERCLIP_DEPLOYMENT_EXPOSURE=private
Environment=PAPERCLIP_PUBLIC_URL=http://localhost:3100
EnvironmentFile=%h/.config/containers/systemd/paperclip.env

[Service]
# Restart=always covers both clean exits (RuntimeMaxSec rolling restart, heap-drain exit 0) and crash exits.
Restart=always
RestartSec=5
TimeoutStartSec=120
# Rolling restart: cap continuous uptime at 6 hours to prevent long-lived OOM accumulation.
RuntimeMaxSec=21600

[Install]
WantedBy=default.target
